# nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'

[program:certbot]
command=/path/to/certbot --nginx -d astraportal.dev-demo.online -d www.astraportal.dev-demo.online
